# packaging/debian/systemd/starrocks-be.service
[Unit]
Description=StarRocks Backend (BE)
After=network.target

[Service]
Type=simple
User=starrocks
Group=starrocks
RuntimeDirectory=starrocks
RuntimeDirectoryMode=0750
ExecStart=/usr/lib/starrocks/be/bin/start_be.sh
Restart=on-failure
LimitNOFILE=262144
LimitMEMLOCK=infinity
TimeoutStartSec=300
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
