# packaging/debian/systemd/starrocks-cn.service
[Unit]
Description=StarRocks Compute Node (CN)
After=network.target

[Service]
Type=simple
User=starrocks
Group=starrocks
RuntimeDirectory=starrocks
ExecStart=/usr/lib/starrocks/cn/bin/start_cn.sh --daemonize 0
Restart=always
RestartSec=5
LimitNOFILE=262144
LimitMEMLOCK=infinity
TimeoutStartSec=300
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
