# packaging/debian/systemd/starrocks-fe.service
[Unit]
Description=StarRocks Frontend (FE)
After=network.target

[Service]
Type=simple
User=starrocks
Group=starrocks
RuntimeDirectory=starrocks
RuntimeDirectoryMode=0750
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -d -m 0750 -o starrocks -g starrocks /run/starrocks
Environment="STARROCKS_HOME=/usr/lib/starrocks/fe"
ExecStart=/usr/lib/starrocks/fe/bin/start_fe.sh
Restart=on-failure
LimitNOFILE=65536
TimeoutStartSec=300
TimeoutStopSec=300

# NOTE:  
# - FE is started via start_fe.sh, which runs in foreground mode by default  
#   (RUN_DAEMON=0) unless the --daemon flag is explicitly provided.  
# - For systemd Type=simple, the main process must stay in the foreground.  
# - Therefore we intentionally do NOT pass --daemon or a separate  
#   "foreground" flag here. BE/CN services are configured in their own  
#   systemd units and startup scripts to keep the main process in the  
#   foreground as well.  

[Install]
WantedBy=multi-user.target
