# packaging/debian/systemd/starrocks-fe.service
[Unit]
Description=StarRocks Frontend (FE)
After=network.target

[Service]
Type=simple
User=starrocks
Group=starrocks
Environment="STARROCKS_HOME=/usr/lib/starrocks/fe"
Environment="PID_DIR=/var/run/starrocks"
# We force daemonize 0 so systemd can track the process correctly
ExecStart=/usr/lib/starrocks/fe/bin/start_fe.sh --conf_dir /etc/starrocks/fe --daemonize 0
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
